# approval-server.service
[Unit]
Description=Snarling Approval Server
After=network.target

[Service]
Type=simple
User=openpi
WorkingDirectory=/home/openpi/.openclaw/workspace/snarling
ExecStart=/home/openpi/displayhat-env/bin/python /home/openpi/.openclaw/workspace/snarling/approval_server.py
Restart=on-failure
RestartSec=5
Environment=OPENCLAW_TOKEN=c1e2798a58fcf2414a4602f743a193838f6e4416eb5a61ed

[Install]
WantedBy=default.target